# contrib/shuck.service
[Unit]
Description=Shuck microVM manager daemon
After=network.target
Documentation=https://github.com/rvben/shuck

[Service]
Type=simple
ExecStart=/usr/local/bin/shuck daemon
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=no
ProtectSystem=strict
ReadWritePaths=/var/lib/shuck
RuntimeDirectory=shuck

[Install]
WantedBy=multi-user.target
