# guest/inittab
# /etc/inittab

::sysinit:/bin/mkdir -p /proc /sys /dev /tmp /var/log /run /lib/modules
::sysinit:/bin/mount -t proc proc /proc
::sysinit:/bin/mount -t sysfs sys /sys
::sysinit:/bin/mount -t devtmpfs devtmpfs /dev
::sysinit:/bin/mkdir -p /dev/pts
::sysinit:/bin/mount -t devpts devpts /dev/pts

# Load vsock modules
::sysinit:/sbin/insmod /lib/modules/vsock.ko
::sysinit:/sbin/insmod /lib/modules/vmw_vsock_virtio_transport_common.ko
::sysinit:/sbin/insmod /lib/modules/vmw_vsock_virtio_transport.ko

# Load network modules
::sysinit:/sbin/insmod /lib/modules/af_packet.ko
::sysinit:/sbin/insmod /lib/modules/failover.ko
::sysinit:/sbin/insmod /lib/modules/net_failover.ko
::sysinit:/sbin/insmod /lib/modules/virtio_net.ko

# Network: bring up loopback and DHCP on eth0
::sysinit:/sbin/ip link set lo up
::sysinit:/sbin/ip link set eth0 up 2>/dev/null
::sysinit:/sbin/udhcpc -b -i eth0

# Start shuck guest agent
::respawn:/usr/local/bin/shuck-agent

# Virtio console for Apple VZ — guarded because Firecracker guests have no
# /dev/hvc0 and respawn would otherwise flood the serial log with open errors.
hvc0::respawn:/bin/sh -c '[ -c /dev/hvc0 ] && exec /sbin/getty -L 115200 hvc0 vt100 || exec sleep 86400'

::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
